# scripts/watcher.service.tmpl
[Unit]
Description=Watcher — personal web monitoring service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory={working_dir}
ExecStart={python_bin} -m watcher.main
Restart=on-failure
RestartSec=10
EnvironmentFile={config_dir}/.env

# Limit runaway resource use
MemoryMax=512M

# Logs go to /tmp/watcher.log (via Python logging in main.py)
StandardOutput=null
StandardError=null

[Install]
WantedBy=default.target
